# ra6m5_cloud_mobile_phone_link_20260703/applications/aliyun_cloud/tcm-cloud.service
[Unit]
Description=TCM Cloud Flask Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tcm-cloud
Environment=APP_HOST=0.0.0.0
Environment=APP_PORT=5000
Environment=BASE_URL=http://YOUR_PUBLIC_IP:5000
ExecStart=/root/tcm-cloud/.venv/bin/python /root/tcm-cloud/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
